# launchd/com.preedee.preed-ee-rebuild.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.preedee.preed-ee-rebuild</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/preedee/Desktop/Cowork/preed-ee &amp;&amp; /Users/preedee/.bun/bin/bun run scan &amp;&amp; git add docs/projects.json docs/index.html &amp;&amp; (git diff --cached --quiet || (git commit -m "daily rescan $(date -u +%Y-%m-%dT%H:%MZ)" &amp;&amp; git push))</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/preedee/Desktop/Cowork/preed-ee/launchd/rebuild.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/preedee/Desktop/Cowork/preed-ee/launchd/rebuild.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/preedee/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
